# conf/supervisord.conf
[supervisord]
nodaemon=true
; PID1 以 root 运行(以便打开 /dev/stdout 写日志)；实际应用进程经各 program 的 user=appuser 降权
user=root

[program:flask]
command=python /app/api.py
directory=/app
user=appuser
environment=HOME="/home/appuser",USER="appuser"
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_maxbytes=0
stderr_maxbytes=0
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:streamlit]
command=streamlit run /app/ui.py --server.port=5002 --server.address=0.0.0.0 --server.headless true
directory=/app
user=appuser
environment=HOME="/home/appuser",USER="appuser"
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_maxbytes=0
stderr_maxbytes=0
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0